# com.ivern.claudeNotice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ivern.claudeNotice</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jeongyuhyeon/.nvm/versions/node/v22.22.2/bin/node</string>
    <string>/Users/jeongyuhyeon/claude-notice/bot.mjs</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <!-- 비정상 종료 시에만 재시작 (설정 파일 없음 등으로 exit 1이면 루프 방지 위해 launchd가 지연 재시작) -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd 기본 PATH에는 nvm 경로가 없으므로 claude CLI 실행을 위해 명시 -->
    <key>PATH</key>
    <string>/Users/jeongyuhyeon/.nvm/versions/node/v22.22.2/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/jeongyuhyeon/.claude/claude-notice/daemon.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jeongyuhyeon/.claude/claude-notice/daemon.log</string>
</dict>
</plist>
